# deploy/hubclock-backend-prod.service
[Unit]
Description=HubClock Backend (Production)
After=network.target mysql.service mariadb.service

[Service]
Type=simple
WorkingDirectory={{PROJECT_ROOT}}
Environment="HUBCLOCK_ENV=production"
ExecStart={{PROJECT_ROOT}}/scripts/run_backend_service.sh --production
Restart=always
RestartSec=5
User={{USER}}

[Install]
WantedBy=multi-user.target
